# branches/dot_experiment/src/dot/install-src/mac/pkg/DotApp/component/Dojo Offline Toolkit/org.dojo.dot.DojoOfflineLaunchd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
	<dict>
		<key>Label</key>
		<string>org.dojo.dot.DojoOfflineLaunchd</string>
		<key>ServiceDescription</key>
		<string>The Dojo Offline user agent, running in the background to provide offline ability to web applications</string>
		<key>ProgramArguments</key>
		<array>
			<string>/Applications/Dojo Offline Toolkit/dotlauncher.sh</string>
			<string>start</string>
		</array>
		<key>OnDemand</key>
		<true/>
		<key>RunAtLoad</key>
		<true/>
		<key>RootDirectory</key>
		<string>/Applications/Dojo Offline Toolkit</string>
		<key>WorkingDirectory</key>
		<string>/Applications/Dojo Offline Toolkit</string>
		<key>StandardOutPath</key>
		<string>/Applications/Dojo Offline Toolkit/dot.out</string>
		<key>StandardErrorPath</key>
		<string>/Applications/Dojo Offline Toolkit/dot.err</string>
	</dict>
</plist>
